# buildroot/package/sysstat/sysstat-10.0.0/cron/sysstat.cron.hourly
#!/bin/sh
# Run system activity accounting tool every 10 minutes
/home/jwang/sourcecode/Perforce/workspace/jxwang_linux-fc8/dev/ruby/soc/branches/main4_adm/buildroot/build_arc/root/lib/sa/sa1 600 6 &
